# raspi/kiosk-pi4.service
[Unit]
Description=Chromium Kiosk (Pi 4/5) - JC Vision Play
After=systemd-user-sessions.service display-manager.service
Wants=graphical.target

[Service]
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Type=simple
Restart=always
RestartSec=3

# Pi 4/5: habilitar aceleração de vídeo e EGL
ExecStart=/usr/bin/chromium-browser \
  --kiosk \
  --noerrdialogs \
  --disable-session-crashed-bubble \
  --incognito \
  --autoplay-policy=no-user-gesture-required \
  --user-data-dir=/home/pi/.kiosk-profile \
  --enable-accelerated-video-decode \
  --ignore-gpu-blocklist \
  --use-gl=egl \
  --start-maximized \
  "https://jc-vision-play.vercel.app/player/c5dc8d9fef4bdf859b1a887e566f5c89"

[Install]
WantedBy=graphical.target